# file-browser.service
[Unit]
Description=File Browser Service
After=network.target

[Service]
# 工作目录（可选）
WorkingDirectory=/opt/file-browser

# 启动前先杀掉可能存在的旧进程
ExecStartPre=/bin/bash -c 'pkill -f "file-browser runserver" || true'

# 启动命令
ExecStart=/opt/file-browser/./file-browser runserver --noreload

# 停止时杀掉进程
ExecStop=/bin/bash -c 'pkill -f "file-browser runserver"'

# 用户和组（可选，根据需要设置）
User=root
Group=root

# 环境变量（如果需要）
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# 进程守护
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target

